# oskari.conf
description "node.js server"

start on started mountall
stop on shutdown

# Automatically respawn:
respawn
respawn limit 99 5

script
    export HOME="/home/jenkins"

    exec /usr/bin/node /srv/oskari.org/app.js >> /var/log/node/oskari.log 2>&1
end script
